# archlive/airootfs/etc/xdg/autostart/konsole.desktop
[Desktop Entry]
Type=Application
Exec=konsole
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Konsole
Comment=Start a terminal
