# archlive/airootfs/etc/xdg/autostart/terminal.desktop
# airootfs/etc/xdg/autostart/terminal.desktop
[Desktop Entry]
Type=Application
Name=Kitty Terminal
Exec=kitty
Terminal=false
Categories=System;
